# quote-system.service
[Unit]
Description=报价管理系统 Quote System
After=network.target

[Service]
Type=simple
User=tong
WorkingDirectory=/opt/quote-system
ExecStart=/home/tong/.local/bin/gunicorn -w 2 -b 127.0.0.1:5000 app:app --access-logfile /opt/quote-system/gunicorn-access.log --error-logfile /opt/quote-system/gunicorn-error.log
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
